# config/data-wallet-demos.service
[Unit]
Description=data-wallet-demos

[Service]
WorkingDirectory=/opt/data-wallet-demos
ExecStart=/usr/bin/npm start &
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=data-wallet-demos
Environment="PORT=80"
Environment="PORTSSL=443"

User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
